# blink_clip_downloader/rootfs/etc/s6-overlay/s6-rc.d/blink-downloader/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Blink Clip Downloader – S6-overlay v3 longrun service.
#
# `with-contenv` loads the HA supervisor environment (SUPERVISOR_TOKEN, etc.)
# before executing this script.  bashio provides structured logging.
#
# As a longrun service, this script must run indefinitely.  S6 will restart it
# automatically if it exits unexpectedly.  A clean shutdown (SIGTERM from S6)
# is handled by the Python app's signal handler.
# ==============================================================================
# shellcheck shell=bash
set -e

bashio::log.info "Starting Blink Clip Downloader..."
bashio::log.info "  Download path : $(bashio::config 'download_path')"
bashio::log.info "  Poll interval : $(bashio::config 'poll_interval') seconds"

exec python3 -m blink_downloader
