# blink_clip_downloader/rootfs/etc/s6-overlay/s6-rc.d/blink-downloader/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Blink Clip Downloader – S6-overlay v3 longrun service entry point.
#
# Registered in s6-rc.d/user/contents.d/blink-downloader (the canonical
# approach used by all official hassio-addons).  with-contenv loads the HA
# supervisor environment (SUPERVISOR_TOKEN etc.) before exec.
# ==============================================================================
set -e

exec python3 -m blink_downloader
