# blink_clip_downloader/rootfs/etc/services.d/blink-downloader/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e
exec python3 -m blink_downloader
